# systemd/dot-config/systemd/user/blue-light-sh.service
[Unit]
Description=Changes blue light filter levels during the day
After=network.target

[Service]
Type=oneshot
ExecStart=notify-send --urgency=normal "Blue light script" "this is set by the blue_light.sh script"
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "wlroots:sway:Wayfire:labwc:Hyprland:niri" ""

[Install]
WantedBy=multi-user.target
